# x.conf
[program:x]
command=python3 /home/curator2/x/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/x/test.err.log
stdout_logfile=/home/curator2/x/test.out.log